# daemon_final.timer
[Unit]
Description=Timer for daemon_final.service

[Timer]
OnBootSec=10s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
